# com.claude.voice-menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.voice-menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.13</string>
        <string>/Users/yashwanth.thota/.claude/scripts/voice-menubar.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/voice-menubar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voice-menubar.log</string>
</dict>
</plist>
